# manifest/etc/services.d/minidlna/run
#!/usr/bin/execlineb -P
minidlnad -d -f /etc/minidlna.conf
